# main/service-restart-lightdm/DEBIAN/postinst
#!/bin/bash

systemctl enable service_restart_lightdm.service