# main/service-restart-lightdm/DEBIAN/prerm
#!/bin/bash
systemctl stop service_restart_lightdm.service
systemctl disable service_restart_lightdm.service
