# main/service-aw859-bluetooth/DEBIAN/postinst
#!/bin/bash

systemctl enable aw859-bluetooth.service