# main/service-aw859-bluetooth/DEBIAN/prerm
#!/bin/bash
systemctl stop aw859-bluetooth.service
systemctl disable aw859-bluetooth.service
