# main/service-disable-fbcon/DEBIAN/postinst
#!/bin/bash

systemctl enable fbcon_dis.service