# main/service-disable-fbcon/DEBIAN/prerm
#!/bin/bash
systemctl stop fbcon_dis.service
systemctl disable fbcon_dis.service
