# main/service-h616-conf-thermal-trip/lib/systemd/system/h616_conf_thermal_trip.service
[Unit]
Description=modify the zone trip
After=basic.target

[Service]
Type=oneshot
ExecStart=/usr/lib/walnutpi/h616_conf_thermal_trip
RemainAfterExit=yes
StandardOutput=null

[Install]
WantedBy=multi-user.target
