# main/service-aw-eth-mac/DEBIAN/postinst
#!/bin/bash

systemctl enable aw_eth_mac.service