# main/service-aw-eth-mac/DEBIAN/prerm
#!/bin/bash
systemctl stop aw_eth_mac.service
systemctl disable aw_eth_mac.service
