# main/service-resize-part/DEBIAN/postinst
#!/bin/bash

systemctl enable resize_part.service