# main/service-h616-mac-detect/DEBIAN/postinst
#!/bin/bash

systemctl enable h616_mac_detect.service