# main/service-h616-mac-detect/DEBIAN/prerm
#!/bin/bash
systemctl stop h616_mac_detect.service
systemctl disable h616_mac_detect.service
