# main/conf-h616-amixer/DEBIAN/postinst
#!/bin/bash

systemctl enable h616_enable_amixer.service