# main/conf-h616-amixer/DEBIAN/prerm
#!/bin/bash
systemctl stop h616_enable_amixer.service
systemctl disable h616_enable_amixer.service
