# main/service-aw869-bluetooth/DEBIAN/postinst
#!/bin/bash

systemctl enable aw869-bluetooth.service