# main/service-aw869-bluetooth/DEBIAN/prerm
#!/bin/bash
systemctl stop aw869-bluetooth.service
systemctl disable aw869-bluetooth.service
