# main/set-lcd/DEBIAN/postinst
#!/bin/bash


cd /usr/lib/walnutpi/set-lcd
chmod +x install
./install

systemctl enable set_lcd_choose.service

/usr/bin/set-lcd hdmi install