# arm64/emmc-bunner/DEBIAN/postinst
#!/bin/bash

set -e

cd /usr/lib/walnutpi/emmc-burner
./install.sh
systemctl enable emmc_burner.service
