# main/service-w2b-dis-eth-led/DEBIAN/postinst
#!/bin/bash

systemctl enable w2b-dis-eth-led.service