# main/service-w2b-dis-eth-led/DEBIAN/prerm
#!/bin/bash
systemctl stop w2b-dis-eth-led.service
systemctl disable w2b-dis-eth-led.service
