# main/service-h616-conf-thermal-trip/DEBIAN/postinst
#!/bin/bash

systemctl enable h616_conf_thermal_trip.service