# main/service-h616-conf-thermal-trip/DEBIAN/prerm
#!/bin/bash
systemctl stop h616_conf_thermal_trip.service
systemctl disable h616_conf_thermal_trip.service
