# main/service-w1b-dis-mmc2/DEBIAN/postinst
#!/bin/bash

systemctl enable w1b_dis_mmc2.service