# main/service-w1b-dis-mmc2/DEBIAN/prerm
#!/bin/bash
systemctl stop w1b_dis_mmc2.service
systemctl disable w1b_dis_mmc2.service
