# main/map-device/DEBIAN/postinst
#!/bin/bash

systemctl enable map_device.service