# main/map-device/DEBIAN/prerm
#!/bin/bash
systemctl stop map_device.service
systemctl disable map_device.service
