# main/hass-status-led/DEBIAN/postinst
#!/bin/bash

systemctl enable hass-status-led.service