# main/hass-status-led/DEBIAN/prerm
#!/bin/bash
systemctl stop hass-status-led.service
systemctl disable hass-status-led.service