# main/conf-t527-audio/DEBIAN/postinst
#!/bin/bash

systemctl enable T527_enable_amixer.service
