# main/conf-t527-audio/DEBIAN/prerm
#!/bin/bash
systemctl stop T527_enable_amixer.service
systemctl disable T527_enable_amixer.service
