# main/goodix-ts-auto-rotation/DEBIAN/postinst
#!/bin/bash

set -e

case "$1" in
    configure)
        systemctl daemon-reload
        systemctl enable goodix-ts-auto-rotation.service
        ;;
esac

exit 0