# main/service-boot-script/DEBIAN/postinst
#!/bin/bash

systemctl enable service_boot_script.service

if [ ! -d /boot/start ]; then
    mkdir -p /boot/start
fi