# main/service-boot-script/DEBIAN/prerm
#!/bin/bash
systemctl stop boot_script.service
systemctl disable boot_script.service
