# main/powerdown-key/DEBIAN/postinst
#!/bin/bash

systemctl enable powerdown-key.service
systemctl start powerdown-key.service