# main/powerdown-key/DEBIAN/prerm
#!/bin/bash
systemctl stop powerdown-key.service
systemctl disable powerdown-key.service